# services/snapraid/run
#!/bin/bash

exec /sbin/setuser abc /usr/local/bin/snapraid --conf="/etc/snapraid/snapraid.conf" sync
